# os-build/systemd/sozu.service.in
[Unit]
Description=Sozu - A HTTP reverse proxy, configurable at runtime, fast and safe, built in Rust.
Documentation=https://docs.rs/sozu/
After=network-online.target
Wants=network-online.target

[Service]
RuntimeDirectory=sozu
PIDFile=__RUNDIR__/sozu/sozu.pid
ExecStart=__BINDIR__/sozu start --config __SYSCONFDIR__/sozu/config.toml
ExecReload=__BINDIR__/sozuctl --config __SYSCONFDIR__/sozu/config.toml reload
Restart=on-failure

[Install]
WantedBy=multi-user.target
