# os-build/sozu.service.in
[Unit]
Description=Sozu - A HTTP reverse proxy, configurable at runtime, fast and safe, built in Rust.
Documentation=https://docs.rs/sozu/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
PIDFile=__STATE_DIR__/sozu.pid
ExecStart=__BIN_DIR__/sozu start --config /etc/sozu/sozu.toml

[Install]
WantedBy=multi-user.target